# Source Code/meta-gpiotest/recipes-gpiotest/gpio-driver/files/gpio-driver.service
[Unit]
Description=GPIO Driver Service
Before=gpio-test.service

[Service]
ExecStart=/sbin/modprobe gpio-driver
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

